# backend/app/mock/attack_coverage.json
[
  { "techniqueID": "T1059.001", "name": "PowerShell", "tactic": "execution", "tested": true, "detected": true, "score": 100, "sigmaRule": "R-21", "comment": "Atomic T1059.001 executed; Sysmon telemetry ingested into Wazuh; Sigma rule R-21 fired." },
  { "techniqueID": "T1003.001", "name": "LSASS Memory", "tactic": "credential-access", "tested": true, "detected": true, "score": 100, "sigmaRule": "R-23", "comment": "Atomic T1003.001 executed; detected by Sigma rule R-23 (comsvcs MiniDump)." },
  { "techniqueID": "T1082", "name": "System Information Discovery", "tactic": "discovery", "tested": true, "detected": false, "score": 40, "sigmaRule": null, "comment": "Atomic T1082 executed; telemetry ingested but no detection rule yet — coverage gap." },
  { "techniqueID": "T1087.002", "name": "Domain Account Discovery", "tactic": "discovery", "tested": true, "detected": false, "score": 40, "sigmaRule": null, "comment": "Executed; SAMR enumeration ingested, detection pending." },
  { "techniqueID": "T1053.005", "name": "Scheduled Task", "tactic": "persistence", "tested": true, "detected": true, "score": 90, "sigmaRule": "R-25", "comment": "Detected via Sigma rule R-25." },
  { "techniqueID": "T1218.011", "name": "Rundll32", "tactic": "defense-evasion", "tested": true, "detected": true, "score": 85, "sigmaRule": null, "comment": "Rundll32 proxy execution detected by Wazuh ruleset." },
  { "techniqueID": "T1046", "name": "Network Service Discovery", "tactic": "discovery", "tested": true, "detected": true, "score": 80, "sigmaRule": null, "comment": "Port scan detected." },
  { "techniqueID": "T1547.001", "name": "Registry Run Keys", "tactic": "persistence", "tested": true, "detected": true, "score": 85, "sigmaRule": null, "comment": "Autostart key creation detected." },
  { "techniqueID": "T1490", "name": "Inhibit System Recovery", "tactic": "impact", "tested": true, "detected": true, "score": 95, "sigmaRule": "R-26", "comment": "vssadmin shadow deletion detected by Sigma rule R-26." },
  { "techniqueID": "T1071.001", "name": "Web Protocols (C2)", "tactic": "command-and-control", "tested": true, "detected": false, "score": 35, "sigmaRule": "R-27", "comment": "Beacon executed; Sigma rule R-27 pending approval — partial coverage." },
  { "techniqueID": "T1558.003", "name": "Kerberoasting", "tactic": "credential-access", "tested": true, "detected": true, "score": 75, "sigmaRule": null, "comment": "RC4 TGS requests detected." },
  { "techniqueID": "T1110.001", "name": "Password Guessing", "tactic": "credential-access", "tested": false, "detected": false, "score": 0, "sigmaRule": null, "comment": "Not yet tested." }
]
